# recipes-core/loadkeys/files/loadkeys.service.j2
[Unit]
Description=Tries to load a certain keyboard layout

[Service]
ExecStart=/bin/sh -c '/usr/bin/loadkeys {{ profile }}'

[Install]
WantedBy=multi-user.target